# infra/copycat-db/copycat-db.timer
[Unit]
Description=Schedule copycat-db

[Timer]
OnCalendar=Daily

[Install]
WantedBy=timers.target
